# config/svtr-sync-scheduler.service
[Unit]
Description=SVTR飞书知识库自动同步调度器
After=network.target

[Service]
Type=simple
User=www-data
WorkingDirectory=/home/lium/chatsvtr
ExecStart=/usr/bin/node scripts/auto-sync-scheduler.js start
ExecStop=/usr/bin/node scripts/auto-sync-scheduler.js stop
Restart=always
RestartSec=10

# 环境变量
Environment=NODE_ENV=production
Environment=PATH=/usr/bin:/usr/local/bin

# 日志
StandardOutput=journal
StandardError=journal
SyslogIdentifier=svtr-sync

# 安全设置
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/home/lium/chatsvtr/logs /home/lium/chatsvtr/assets/data /home/lium/chatsvtr/config

[Install]
WantedBy=multi-user.target